# services/com.corrkit.watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent for corrkit watch
  Install: cp com.corrkit.watch.plist ~/Library/LaunchAgents/
  Load:    launchctl load ~/Library/LaunchAgents/com.corrkit.watch.plist
  Logs:    tail -f /tmp/corrkit-watch.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.corrkit.watch</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>uv</string>
        <string>run</string>
        <string>corrkit</string>
        <string>watch</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/YOU/work/btakita/corrkit</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/corrkit-watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/corrkit-watch.log</string>
</dict>
</plist>
